# debian/postinst
#!/bin/sh

set -e

invoke-rc.d procps restart
